# actividad4/saludo.service
[Unit]
Description=Saludo y fecha actual

[Service]
Type=simple
ExecStart=/home/droid/Desktop/saludo.sh

[Install]
WantedBy=multi-user.target